# conf.d/system/udev.d/95-external-head-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/octocat/.Xauthority", RUN+="xrandr --output VGA1 --auto --right-of LVDS1"