# deploy/roles/services/templates/job-vpn.openrc
#!/sbin/openrc-run

VPNCONF="{{ custom_vpn_datadir }}/office.ovpn"
PIDFILE={{ pids_dir }}/job-vpn.pid
DAEMON_OPTS="--config ${VPNCONF} --writepid ${PIDFILE} --daemon"

start () {
	ebegin "Starting job vpn"
	start-stop-daemon --background --start --exec /usr/sbin/openvpn --make-pidfile --pidfile ${PIDFILE} \
             -- ${DAEMON_OPTS}
	eend $? "Failed to start job vpn"
}

stop() {
	ebegin "Stopping job vpn"
	start-stop-daemon --stop --pidfile ${PIDFILE}
	eend $? "Failed to stop job vpn"
}

