# deploy/roles/services/templates/shepherd.openrc
#!/sbin/openrc-run

DAEMON=/usr/local/bin/shepherd
PIDFILE={{ pids_dir }}/shepherd.pid

start () {
	ebegin "Starting shepherd"
	start-stop-daemon --background --start -u octocat:octocat --exec ${DAEMON} \
        --make-pidfile --pidfile ${PIDFILE}
	eend $? "Failed to start shepherd"
}

stop() {
	ebegin "Stopping shepherd"
	start-stop-daemon --stop --pidfile ${PIDFILE}
	eend $? "Failed to stop shepherd"
}
